# config/supervisord.conf
[unix_http_server]
file=/var/run/supervisord.sock
chmod=0700

[supervisord]
logfile=/home/androidusr/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisord.sock

[program:emulator]
command=/usr/local/bin/start-emulator.sh
user=androidusr
autostart=true
autorestart=true
startretries=3
startsecs=10
stopwaitsecs=30
priority=1
stdout_logfile=/home/androidusr/logs/emulator.log
stderr_logfile=/home/androidusr/logs/emulator-err.log

[program:ws-scrcpy]
command=/usr/local/bin/start-ws-scrcpy.sh
user=androidusr
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
priority=2
stdout_logfile=/home/androidusr/logs/ws-scrcpy.log
stderr_logfile=/home/androidusr/logs/ws-scrcpy-err.log

[program:nport]
command=/usr/local/bin/start-tunnel.sh
user=androidusr
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
priority=3
stdout_logfile=/home/androidusr/logs/nport.log
stderr_logfile=/home/androidusr/logs/nport-err.log

[program:healthcheck]
command=/usr/local/bin/health-check-loop.sh
user=androidusr
autostart=true
autorestart=true
startretries=2
startsecs=120
priority=10
stdout_logfile=/home/androidusr/logs/health.log

[program:session-monitor]
command=/usr/local/bin/session-monitor.sh
user=androidusr
autostart=true
autorestart=false
startretries=1
startsecs=5
priority=11
stdout_logfile=/home/androidusr/logs/session.log
